# roles/router/fabio/config/fabio.service
[Unit]
Description=Fabio

[Service]
Restart=always
ExecStart=/opt/fabio/bin/fabio -cfg /opt/fabio/config.properties

[Install]
WantedBy=local.target